# rootfs/etc/s6-overlay/s6-rc.d/headscale/run
#!/command/with-contenv sh
exec "${BASE_PATH}/headscale/headscale" "serve"